# deployment/desktop-systemd-user/claude-session-atlas-refresh@.service
[Unit]
Description=Claude session atlas refresh handler (per-connection, socket-activated)

[Service]
ExecStart=/home/kjdragan/lrepos/universal_agent/scripts/claude-session-atlas --refresh-http
StandardInput=socket
StandardOutput=socket
StandardError=journal
TimeoutStartSec=600
