# deployment/desktop-systemd-user/claude-session-atlas.service
[Unit]
Description=Claude session atlas — regenerate + publish the living session map (desktop)

[Service]
Type=oneshot
ExecStart=/bin/bash -lc '/home/kjdragan/.local/bin/claude-session-atlas --publish'
TimeoutStartSec=600
