# deployment/systemd-user/ua-youtube-forward-tunnel.service
[Unit]
Description=UA reverse tunnel (VPS -> local gateway) for YouTube hook mirroring
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Template unit: install via scripts/install_ua_youtube_forward_tunnel_user_service.sh
# which renders __REPO_ROOT__ to your actual checkout path.
WorkingDirectory=__REPO_ROOT__
ExecStart=__REPO_ROOT__/scripts/forward_youtube_hooks_to_local.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
