# .agents/skills/project-scaffolder/templates/infrastructure/api.service.tmpl
[Unit]
Description={{PROJECT_NAME}} API (FastAPI/Uvicorn)
After=network.target {{PROJECT_NAME}}-db.service
Requires={{PROJECT_NAME}}-db.service

[Service]
Type=simple
User=ua
WorkingDirectory={{PROJECT_DIR}}/backend
Environment=PATH={{PROJECT_DIR}}/backend/.venv/bin:/usr/local/bin:/usr/bin
ExecStart={{PROJECT_DIR}}/backend/.venv/bin/uvicorn {{PROJECT_MODULE}}.main:app --host 127.0.0.1 --port {{BACKEND_PORT}} --workers 2
Restart=on-failure
RestartSec=5
MemoryMax=512M
MemoryHigh=384M

[Install]
WantedBy=multi-user.target
