# .agents/skills/project-scaffolder/templates/infrastructure/docs.service.tmpl
[Unit]
Description={{PROJECT_NAME}} Documentation Server (MkDocs)
After=network.target

[Service]
Type=simple
User=ua
WorkingDirectory={{PROJECT_DIR}}
ExecStart={{PROJECT_DIR}}/backend/.venv/bin/python -m http.server {{DOCS_PORT}} --directory site
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
