# CSI_Ingester/development/deployment/systemd/csi-hourly-token-report.timer
[Unit]
Description=Run CSI hourly token usage report at minute 5

[Timer]
OnCalendar=*:05
Persistent=true
Unit=csi-hourly-token-report.service

[Install]
WantedBy=timers.target

