# CSI_Ingester/development/deployment/systemd/csi-reddit-trend-report.timer
[Unit]
Description=Run CSI Reddit trend report 3x daily (7am/1pm/7pm CT)
PartOf=csi.target

[Timer]
OnCalendar=*-*-* 12:18:00 UTC
OnCalendar=*-*-* 18:18:00 UTC
OnCalendar=*-*-* 00:18:00 UTC
Persistent=true

[Install]
WantedBy=csi.target
