# CSI_Ingester/development/deployment/systemd/csi-threads-semantic-enrich.timer
[Unit]
Description=Run CSI Threads semantic enrichment every 15 minutes

[Timer]
OnCalendar=*:11/15
Persistent=true
Unit=csi-threads-semantic-enrich.service

[Install]
WantedBy=timers.target
