# CSI_Ingester/development/deployment/systemd/csi-delivery-slo-gatekeeper.service
[Unit]
Description=Run CSI daily delivery reliability SLO gatekeeper
After=network.target csi-ingester.service universal-agent-gateway.service
Wants=csi-ingester.service

[Service]
Type=oneshot
WorkingDirectory=/opt/universal_agent/CSI_Ingester/development
EnvironmentFile=/opt/universal_agent/CSI_Ingester/development/deployment/systemd/csi-ingester.env
ExecStart=/bin/bash -lc 'DAY="$(date -u -d "yesterday" +%%F)"; /opt/universal_agent/CSI_Ingester/development/scripts/csi_run.sh /opt/universal_agent/CSI_Ingester/development/.venv/bin/python /opt/universal_agent/CSI_Ingester/development/scripts/csi_delivery_slo_gatekeeper.py --db-path /opt/universal_agent/CSI_Ingester/development/var/csi.db --day "$DAY"'

[Install]
WantedBy=multi-user.target
