# CSI_Ingester/development/deployment/systemd/csi-global-trend-brief.timer
[Unit]
Description=Run CSI global trend brief every 2 hours

[Timer]
OnCalendar=*-*-* 0/2:05:00
Persistent=true
Unit=csi-global-trend-brief.service

[Install]
WantedBy=timers.target
