# CSI_Ingester/development/deployment/systemd/csi-threads-webhook-canary-verify.timer
[Unit]
Description=Run Threads webhook canary verification every 2 hours

[Timer]
OnCalendar=*-*-* 00/2:00:00
Persistent=true
Unit=csi-threads-webhook-canary-verify.service

[Install]
WantedBy=timers.target
