# CSI_Ingester/development/deployment/systemd/csi-reddit-telegram-digest.timer
[Unit]
Description=Run CSI Reddit Telegram digest every 10 minutes

[Timer]
OnCalendar=*:1/10
Persistent=true
Unit=csi-reddit-telegram-digest.service

[Install]
WantedBy=timers.target
