# CSI_Ingester/development/deployment/systemd/csi-replay-dlq.timer
[Unit]
Description=Run CSI DLQ replay every 15 minutes

[Timer]
OnCalendar=*:08/15
Persistent=true
Unit=csi-replay-dlq.service

[Install]
WantedBy=timers.target

