# CSI_Ingester/development/deployment/systemd/csi-replay-dlq.service
[Unit]
Description=Replay CSI dead-letter queue entries to UA
After=network.target csi-ingester.service universal-agent-gateway.service
Wants=csi-ingester.service

[Service]
Type=oneshot
WorkingDirectory=/opt/universal_agent/CSI_Ingester/development
EnvironmentFile=/opt/universal_agent/CSI_Ingester/development/deployment/systemd/csi-ingester.env
ExecStart=/opt/universal_agent/CSI_Ingester/development/scripts/csi_run.sh /opt/universal_agent/CSI_Ingester/development/.venv/bin/python /opt/universal_agent/CSI_Ingester/development/scripts/csi_replay_dlq.py --db-path /opt/universal_agent/CSI_Ingester/development/var/csi.db --limit 50 --max-attempts 3

[Install]
WantedBy=multi-user.target

