# CSI_Ingester/development/deployment/systemd/csi-hourly-token-report.service
[Unit]
Description=CSI Hourly Token Usage Report to UA
After=network.target csi-ingester.service
Wants=csi-ingester.service

[Service]
Type=oneshot
WorkingDirectory=/opt/universal_agent/CSI_Ingester/development
EnvironmentFile=/opt/universal_agent/CSI_Ingester/development/deployment/systemd/csi-ingester.env
ExecStart=/opt/universal_agent/CSI_Ingester/development/scripts/csi_run.sh /opt/universal_agent/CSI_Ingester/development/.venv/bin/python /opt/universal_agent/CSI_Ingester/development/scripts/csi_hourly_token_report.py --db-path /opt/universal_agent/CSI_Ingester/development/var/csi.db
