# CSI_Ingester/development/deployment/systemd/csi-category-quality-loop.timer
[Unit]
Description=Run CSI category quality loop hourly

[Timer]
OnCalendar=*:27
Persistent=true
Unit=csi-category-quality-loop.service

[Install]
WantedBy=timers.target
