# CSI_Ingester/development/deployment/systemd/csi-threads-rollout-verify.timer
[Unit]
Description=Run Threads rollout verification daily after token refresh

[Timer]
OnCalendar=*-*-* 03:35:00 UTC
Persistent=true
Unit=csi-threads-rollout-verify.service

[Install]
WantedBy=timers.target

