# CSI_Ingester/development/deployment/systemd/csi-quality-assessment.service
[Unit]
Description=Run CSI source quality assessment daily (evening)
PartOf=csi.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/universal_agent/CSI_Ingester/development
ExecStart=/opt/universal_agent/CSI_Ingester/development/.venv/bin/python scripts/csi_source_quality_assessment.py --lookback-days 7
Environment=CSI_DB_PATH=/opt/universal_agent/CSI_Ingester/development/var/csi.db
TimeoutStartSec=300
