# CSI_Ingester/development/deployment/systemd/csi-quality-assessment.service
[Unit]
Description=Run CSI source quality assessment daily (evening)
PartOf=csi.target

[Service]
Type=oneshot
User=ua
WorkingDirectory=/opt/universal_agent/CSI_Ingester/development
ExecStart=/opt/universal_agent/CSI_Ingester/development/.venv/bin/python scripts/csi_source_quality_assessment.py --lookback-days 7
# INTERIM: deliberately pointed at the (empty) dev DB, NOT the canonical
# /var/lib/universal-agent/csi/csi.db. Pointing it at canonical (PR #1086) made
# every run fail with SQLITE_BUSY — the live csi-ingester is the sole, continuous
# writer and holds the WAL write lock, so an external batch writer never gets a
# window (verified: BEGIN IMMEDIATE failed after 30s). Real fix = move source-
# quality scoring in-process into the ingester (where it already holds the lock);
# tracked as a follow-up. Until then this stays a harmless no-op.
Environment=CSI_DB_PATH=/opt/universal_agent/CSI_Ingester/development/var/csi.db
TimeoutStartSec=300
