# CSI_Ingester/development/deployment/systemd/csi-ingester.service
[Unit]
Description=CSI Ingester Service
After=network.target
PartOf=csi.target

[Service]
Type=simple
WorkingDirectory=/opt/universal_agent/CSI_Ingester/development
# DO NOT wrap ExecStart in csi_run.sh / `infisical run`, and DO NOT add
# EnvironmentFile=-/opt/universal_agent/.env here. Tried in #820/#824 to give
# batch_brief an LLM key, but `infisical run` injects ALL ~289 vault secrets and
# OVERRODE this continuous ingester's CSI-specific proxy creds from
# csi-ingester.env -> every YouTube RSS fetch failed with `error=407 NO_USER` and
# the circuit breaker aborted the ingest cycle (reverted here). The oneshot timer
# units can take the full vault env; this continuous ingester cannot. To give
# batch_brief an LLM key, seed the DEDICATED CSI lane instead (CSI_LLM_AUTH_MODE=1
# + CSI_ZAI_API_KEY/CSI_ANTHROPIC_API_KEY in csi-ingester.env; see
# scripts/csi_seed_dedicated_auth_env.py) — that adds only the key, not the proxy.
EnvironmentFile=/opt/universal_agent/CSI_Ingester/development/deployment/systemd/csi-ingester.env
ExecStart=/usr/bin/env uv run uvicorn csi_ingester.app:app --host 127.0.0.1 --port 8091
MemoryHigh=900M
MemoryMax=1400M
MemorySwapMax=1024M
TasksMax=2048
OOMPolicy=continue
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
