# CSI_Ingester/development/deployment/systemd/csi-delivery-health-canary.service
[Unit]
Description=Run CSI delivery-health runtime canary
After=network.target csi-ingester.service universal-agent-gateway.service
Wants=csi-ingester.service

[Service]
Type=oneshot
WorkingDirectory=/opt/universal_agent/CSI_Ingester/development
EnvironmentFile=/opt/universal_agent/CSI_Ingester/development/deployment/systemd/csi-ingester.env
ExecStart=/opt/universal_agent/CSI_Ingester/development/scripts/csi_run.sh /opt/universal_agent/CSI_Ingester/development/.venv/bin/python /opt/universal_agent/CSI_Ingester/development/scripts/csi_delivery_health_canary.py --db-path /opt/universal_agent/CSI_Ingester/development/var/csi.db

[Install]
WantedBy=multi-user.target

