# CSI_Ingester/development/deployment/systemd/csi-rss-quality-gate.timer
[Unit]
Description=Run CSI RSS quality gate checks every 15 minutes

[Timer]
OnCalendar=*:00/15
Persistent=true
Unit=csi-rss-quality-gate.service

[Install]
WantedBy=timers.target
