# deployment/systemd/universal-agent-service-watchdog.service
[Unit]
Description=Universal Agent Service Watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/universal_agent
EnvironmentFile=-/opt/universal_agent/.env
Environment=UA_WATCHDOG_HEALTH_FAIL_THRESHOLD=10
ExecStart=/opt/universal_agent/scripts/vps_service_watchdog.sh

[Install]
WantedBy=multi-user.target
