# deployment/systemd/universal-agent-youtube-oauth-watchdog.service
[Unit]
Description=Universal Agent YouTube OAuth watchdog (S5 Phase A batch A4 — deploy-independent daily timer, SECRET-BEARING: OAuth refresh token)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ua
Group=ua
WorkingDirectory=/opt/universal_agent
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/universal_agent/src
Environment=PYDANTIC_DISABLE_PLUGINS=logfire-plugin
# Secrets-bootstrap backstop (mirrors universal-agent-proactive-artifact-digest)
# — LOAD BEARING. The watchdog tests the YouTube OAuth refresh token's liveness
# + age and emails a re-auth button via AgentMail. main() calls
# initialize_runtime_secrets() with NO profile arg so it honors
# UA_DEPLOYMENT_PROFILE below; the inline vars are the .env-absent backstop
# forcing the strict `vps` production load (else the token check + AgentMail
# send fail keyless and the operator never gets the warning). .env OVERRIDES
# these when present.
Environment=UA_DEPLOYMENT_PROFILE=vps
Environment=UA_INFISICAL_ENABLED=1
Environment=INFISICAL_ENVIRONMENT=production
EnvironmentFile=-/opt/universal_agent/.env
ExecStart=/opt/universal_agent/.venv/bin/python -m universal_agent.scripts.youtube_oauth_watchdog
MemoryMax=1G
# Pure stdlib + httpx token probe + at most one AgentMail send. Bounded to its
# in-process budget (timeout_seconds=120) so a hung probe self-heals next slot.
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
