# deployment/systemd/universal-agent-nightly-wiki.service
[Unit]
Description=Universal Agent nightly wiki generation (S5 Phase A batch A4 — deploy-independent daily timer, SECRET-BEARING: NotebookLM cookie + LLM)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ua
Group=ua
WorkingDirectory=/opt/universal_agent
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/universal_agent/src
Environment=PYDANTIC_DISABLE_PLUGINS=logfire-plugin
# Secrets-bootstrap backstop (mirrors universal-agent-proactive-artifact-digest)
# — LOAD BEARING. nightly_wiki_agent dispatches a VP LLM mission that drives
# `nlm` against NotebookLM (auth = NOTEBOOKLM_AUTH_COOKIE_HEADER). main() was
# fixed (M2) from initialize_runtime_secrets(profile="local_workstation") to a
# bare call so it honors UA_DEPLOYMENT_PROFILE below; the inline vars are the
# .env-absent backstop forcing the strict `vps` production load (else the
# NotebookLM cookie + LLM keys are absent and the run dies exit 1). .env
# OVERRIDES these when present.
Environment=UA_DEPLOYMENT_PROFILE=vps
Environment=UA_INFISICAL_ENABLED=1
Environment=INFISICAL_ENVIRONMENT=production
EnvironmentFile=-/opt/universal_agent/.env
ExecStart=/opt/universal_agent/.venv/bin/python -m universal_agent.scripts.nightly_wiki_agent
MemoryMax=1G
# Dispatches a VP LLM mission over pending CSI cards + NotebookLM round-trips.
# Bounded to its in-process budget (timeout_seconds=1800) so a hung mission
# self-heals on the next night's slot.
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
