# deployment/systemd/universal-agent-csi-demo-triage-rank.timer
[Unit]
Description=Run the Universal Agent CSI demo-triage LLM ranking at 10:05 and 15:05 America/Chicago (deploy-independent)

[Timer]
# OnCalendar wall-clock anchor — NEVER monotonic. cron `5 10,15 * * *`
# America/Chicago (minute 05 of hours 10 and 15 — twice daily, both inside the
# active window), DST handled by systemd; Persistent replays a missed slot
# inside a deploy window. systemd-255 honors the OnCalendar timezone suffix and
# the hour list `10,15` expands to one fire per listed hour (verified with
# `systemd-analyze calendar`).
OnCalendar=*-*-* 10,15:05:00 America/Chicago
Persistent=true
RandomizedDelaySec=120
AccuracySec=1min
Unit=universal-agent-csi-demo-triage-rank.service

[Install]
WantedBy=timers.target
