# deployment/systemd/universal-agent-proactive-signal-card-sync.service
[Unit]
Description=Universal Agent proactive signal card generation tick (autonomous hourly card sync)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ua
Group=ua
WorkingDirectory=/opt/universal_agent
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/universal_agent/src
Environment=PYDANTIC_DISABLE_PLUGINS=logfire-plugin
# Deployment-profile backstop (mirrors the batch-1 pure-SQLite units). .env
# (loaded below) OVERRIDES these inline vars when present; they only matter if
# .env is absent. This tick is pure SQLite — no LLM, no mail, no secrets — it
# reads the CSI/Discord feedstock and writes proactive_signal_cards in
# activity_state.db; the inline vars only ensure consistent path/profile
# resolution if .env were ever missing.
Environment=UA_DEPLOYMENT_PROFILE=vps
Environment=UA_INFISICAL_ENABLED=1
Environment=INFISICAL_ENVIRONMENT=production
EnvironmentFile=-/opt/universal_agent/.env
# Generates/refreshes YouTube diamond + Discord signal cards and runs the TTL
# sweep (proactive_signals.generate_signal_cards). Exits 0 on success (incl.
# dormant no-op), 1 on a generation error. Pure SQLite, so TimeoutStartSec is
# left at the oneshot default — but bounded modestly since it reads a large CSI
# DB.
TimeoutStartSec=600
ExecStart=/opt/universal_agent/.venv/bin/python -m universal_agent.scripts.proactive_signal_card_sync
MemoryMax=1G

[Install]
WantedBy=multi-user.target
