# deployment/systemd/universal-agent-proactive-report-afternoon.timer
[Unit]
Description=Fire the Universal Agent afternoon proactive report at 16:05 America/Chicago (deploy-independent)

[Timer]
# OnCalendar wall-clock anchor — NEVER monotonic. cron `5 16 * * *`
# America/Chicago, DST handled by systemd; Persistent replays a missed slot.
# Shares universal-agent-proactive-report.service with the morning/midday
# timers.
OnCalendar=*-*-* 16:05:00 America/Chicago
Persistent=true
RandomizedDelaySec=120
AccuracySec=1min
Unit=universal-agent-proactive-report.service

[Install]
WantedBy=timers.target
