# deployment/systemd/universal-agent-csi-convergence-sync.timer
[Unit]
Description=Run the Universal Agent CSI convergence sync at the top of every hour (00..23) America/Chicago; runtime-gated to the active window unless UA_CSI_CONVERGENCE_SYNC_24_7 is set (deploy-independent)

[Timer]
# OnCalendar wall-clock anchor — NEVER monotonic. Fires minute 0 of every hour
# (00..23) America/Chicago, DST handled by systemd; Persistent replays a missed
# slot inside a deploy window. systemd-255 honors the OnCalendar timezone suffix
# and the H..H hour range expands to one fire per hour (verified with
# `systemd-analyze calendar`).
#
# 24/7 SCHEDULE, RUNTIME-GATED: the schedule is full-day so the per-run dormancy
# decision lives in the ExecStart script (csi_convergence_sync.py), which calls
# gates via dormancy.should_run, controlled by UA_CSI_CONVERGENCE_SYNC_24_7.
# Default stays windowed (work no-ops outside 6 AM–10 PM Houston); set
# UA_CSI_CONVERGENCE_SYNC_24_7=true in Infisical to run 24/7 — no timer edit /
# reinstall needed to flip.
OnCalendar=*-*-* 00..23:00:00 America/Chicago
Persistent=true
RandomizedDelaySec=120
AccuracySec=1min
Unit=universal-agent-csi-convergence-sync.service

[Install]
WantedBy=timers.target
