# deployment/systemd/universal-agent-hourly-intel-digest.service
[Unit]
Description=Universal Agent deterministic hourly intel-digest delivery (S5 Phase A batch 3 — deploy-independent hourly timer)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ua
Group=ua
WorkingDirectory=/opt/universal_agent
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/universal_agent/src
Environment=PYDANTIC_DISABLE_PLUGINS=logfire-plugin
# Secrets-bootstrap backstop (mirrors universal-agent-proactive-artifact-digest)
# — LOAD BEARING. hourly_intel_digest_cron sends an AgentMail convergence-brief
# digest + reads the activity DB. Its main() now calls initialize_runtime_secrets()
# with NO profile arg so it honors UA_DEPLOYMENT_PROFILE below; the inline vars
# are the .env-absent backstop forcing the strict `vps` production load (else a
# keyless run). .env OVERRIDES these when present.
Environment=UA_DEPLOYMENT_PROFILE=vps
Environment=UA_INFISICAL_ENABLED=1
Environment=INFISICAL_ENVIRONMENT=production
EnvironmentFile=-/opt/universal_agent/.env
# Exits 0 on a throttled/no-op hour; exit 1 reserved for genuine errors.
ExecStart=/opt/universal_agent/.venv/bin/python -m universal_agent.scripts.hourly_intel_digest_cron
MemoryMax=1G
# AgentMail WebSocket can hang; bound to the job's old in-process budget
# (timeout_seconds=300) so a hung run self-heals on the next slot.
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
