# deployment/systemd/universal-agent-architecture-canvas-drift.service
[Unit]
Description=Universal Agent Architecture Canvas pointer drift check (S5 Phase A batch 1 — deploy-independent weekly timer)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ua
Group=ua
WorkingDirectory=/opt/universal_agent
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/opt/universal_agent/src
Environment=PYDANTIC_DISABLE_PLUGINS=logfire-plugin
# Secrets-bootstrap backstop (mirrors universal-agent-proactive-health). .env
# OVERRIDES these when present. architecture_canvas_drift_check is pure FS/git
# (globs docs/architecture-view/sources/*.yaml, shells `git log`, writes a local
# markdown artifact) with NO Infisical-resolved secrets — these vars are inert
# for this job and kept only for unit uniformity.
Environment=UA_DEPLOYMENT_PROFILE=vps
Environment=UA_INFISICAL_ENABLED=1
Environment=INFISICAL_ENVIRONMENT=production
EnvironmentFile=-/opt/universal_agent/.env
# By design this returns exit 1 when a `source:` pointer path is missing (drift)
# and exit 2 on internal error, so the oneshot .service will enter `failed`
# state on a genuine finding — this MIRRORS the prior "failed cron tick" signal
# and is the intended surface, NOT a crash. It exits 0 on clean / stale-only
# (writes the drift report then exits 0). A oneshot is deliberately NOT in the
# is-active vps_service_watchdog specs, so a `failed` run does not false-page.
ExecStart=/opt/universal_agent/.venv/bin/python -m universal_agent.scripts.architecture_canvas_drift_check
MemoryMax=1G

[Install]
WantedBy=multi-user.target
