# discord_intelligence/ua-discord-intelligence.service
[Unit]
Description=Universal Agent - Discord Intelligence Daemon
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/universal_agent
# Loads .env or infisical automatically based on UA standards
Environment="PYTHONPATH=/opt/universal_agent/src:/opt/universal_agent"
ExecStart=/usr/local/bin/uv run python -m discord_intelligence.daemon
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
