# deployment/systemd/templates/ua-discord-intelligence.service.template
[Unit]
Description=Universal Agent - Discord Intelligence Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_ROOT__
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=__APP_ROOT__/src:__APP_ROOT__
Environment=PYDANTIC_DISABLE_PLUGINS=logfire-plugin
Environment=UA_DISCORD_AUTO_CREATE_RELEASE_TASKS=0
EnvironmentFile=-__APP_ROOT__/.env
ExecStart=__APP_ROOT__/.venv/bin/python -m discord_intelligence.daemon
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
