# deployment/systemd/templates/universal-agent-gateway.service.template
[Unit]
Description=__UNIT_DESCRIPTION__
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
Group=__APP_USER__
WorkingDirectory=__APP_ROOT__
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=__APP_ROOT__/src
Environment=PYDANTIC_DISABLE_PLUGINS=logfire-plugin
EnvironmentFile=-__APP_ROOT__/.env
ExecStart=__APP_ROOT__/.venv/bin/python -m universal_agent.gateway_server
Restart=always
RestartSec=5
# Graceful-shutdown window (ADR §C2 — deploy-restart resilience). systemd waits up
# to TimeoutStopSec after SIGTERM before SIGKILL; the gateway uses this window to
# DRAIN an in-flight Simone heartbeat iteration (UA_HEARTBEAT_DRAIN_TIMEOUT_SECONDS,
# default 45s) instead of killing it mid-flight (harm H2). 90s == systemd's current
# default — pinned here so the drain budget always keeps headroom if that default moves.
TimeoutStopSec=90s
# Resource limits — prevent process explosion from consuming all VPS memory.
# MemoryMax: hard cgroup limit (systemd kills cgroup if exceeded).
# MemoryHigh: soft limit (triggers kernel memory reclaim/swap pressure).
# TasksMax: caps total process count within the cgroup.
MemoryMax=8G
MemoryHigh=6G
TasksMax=500

[Install]
WantedBy=multi-user.target
